# scripts/systemd/iptv-panel-autostart.service
[Unit]
Description=IPTV Panel - Autostart channels after reboot
Wants=redis-server.service
After=network.target redis-server.service

[Service]
Type=oneshot
WorkingDirectory=__IPTV_PANEL_DIR__
User=www-data
Group=www-data

# IMPORTANT:
# The command intentionally spawns long-running ffmpeg processes detached.
# Prevent systemd from killing those child processes when the oneshot exits.
KillMode=none
TimeoutStartSec=60

# Starts channels that are marked as status=live in DB but not currently running.
ExecStart=/usr/bin/php8.4 artisan channels:autostart --only-missing

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
# NOTE: ProtectHome=true can block access when the project is under a user home directory.
ProtectHome=false
ReadWritePaths=__IPTV_PANEL_DIR__/storage __IPTV_PANEL_DIR__/bootstrap/cache

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
