# scripts/systemd/iptv-panel-schedule.service
[Unit]
Description=IPTV Panel - Laravel Scheduler (schedule:work)
Wants=redis-server.service
After=network.target redis-server.service

[Service]
Type=simple
WorkingDirectory=__IPTV_PANEL_DIR__
# IMPORTANT: change User/Group if your deployment uses a different account
User=www-data
Group=www-data

ExecStart=/usr/bin/php8.4 artisan schedule:work
Restart=always
RestartSec=2

# Optional hardening (safe defaults)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
# NOTE: ProtectHome=true can block access when the project is under a user home directory.
ProtectHome=false
ReadWritePaths=__IPTV_PANEL_DIR__/storage __IPTV_PANEL_DIR__/bootstrap/cache

# Logs go to journald (view with: journalctl -u iptv-panel-schedule -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
